# systemd/mnt-smb-media.automount
[Unit]
Description=samba automount for media 

[Automount]
Where=/mnt/smb/media
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
