# systemd/mnt-smb-books.automount
[Unit]
Description=samba automount for books 

[Automount]
Where=/mnt/smb/books
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
